# streamrecord.conf
[program:streamrecord]
command=/usr/bin/uwsgi --ini /etc/uwsgi/apps-enabled/streamrecord.ini --uid streamrecord --gid streamrecord
autorestart=true

